# cache.target
[Unit]
Description=Cache target

[Install]
WantedBy=default.target
