# kde@.target
[Unit]
Description=KDE target
Requires=default.target
